# packaging/systemd/dpm-node.service
[Unit]
Description=DPM Node Service
After=network-online.target
Wants=network-online.target

[Service]
# Run as a dedicated unprivileged user created by the installer
User=dpmuser
Group=dpmuser
Type=simple

# Working directory for the node implementation (installed under /opt/dpm)
WorkingDirectory=/opt/dpm/node

# Ensure Python finds the installed packages under /opt/dpm
Environment=PYTHONPATH=/opt/dpm

# Command to execute the node (uses env to pick the active python3)
ExecStart=/usr/bin/env python3 /opt/dpm/node/node.py

# Restart behavior
Restart=on-failure
RestartSec=5

# Allow adjusting process priority without running as root
AmbientCapabilities=CAP_SYS_NICE
CapabilityBoundingSet=CAP_SYS_NICE

# Log to journal
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
