# debian/dpmd.prerm
#!/bin/sh
set -e

case "$1" in
    remove|upgrade|deconfigure)
        if [ -d /run/systemd/system ]; then
            systemctl stop dpmd.service || true
            systemctl disable dpmd.service || true
        fi
        ;;
esac

#DEBHELPER#

exit 0
