# agent/dpm-agent.service
# filepath: /etc/systemd/system/dpm-agent.service
[Unit]
Description=DPM Agent Service
After=network-online.target # Ensure network is up for LCM
Wants=network-online.target

[Service]
# Option 1: Run as a dedicated user (replace 'dpmuser' with an actual user)
# User=dpmuser
# Group=dpmuser

# Option 2: Run as root (simpler, but less secure, needed for real-time priority without extra config)
User=root

# Set the working directory to where agent.py is located
# IMPORTANT: Replace /home/mbustos/dpm/agent with the actual path
WorkingDirectory=/home/mbustos/dpm/agent

# Command to execute the agent
# IMPORTANT: Ensure python3 points to the correct interpreter if using virtual environments
ExecStart=/usr/bin/python3 /home/mbustos/dpm/agent/agent.py

# Restart the service if it fails
Restart=on-failure
RestartSec=5s

# Standard output and error logging will go to journald
StandardOutput=journal
StandardError=journal

# If running as non-root and needing real-time priority:
# AmbientCapabilities=CAP_SYS_NICE

[Install]
WantedBy=multi-user.target